# first-boot/etc/systemd/system/first-boot.service
[Unit]
Description=Post imaging script

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/first-boot.sh

[Install]
WantedBy=multi-user.target
